# examples/systemd/docker-ztd-restore.service
[Unit]
Description=Restore docker-ztd registered projects after Docker startup
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Set the account that has access to Docker and docker-ztd plugin.
User=root
# Optional: set if your plugin or state depends on a specific home directory.
Environment=HOME=/root
# Wait until Docker daemon is responsive.
ExecStartPre=/usr/bin/docker info
ExecStart=/usr/bin/docker ztd restore --all
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
